# arch_package/systemd.service
[Unit]
Description=Tarock server
After=network.target

[Service]
ExecStart=/usr/bin/tarock-server
User=tarock

[Install]
WantedBy=multi-user.target
