# arch_package/tarock-server.service
[Unit]
Description=Tarock server
After=network.target

[Service]
ExecStart=/usr/bin/tarock-server

[Install]
WantedBy=multi-user.target
